# tokogo.service
[Unit]
Description=TokoGo API Server
After=network.target mysql.service
Requires=mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/tokogo
ExecStart=/root/tokogo/main
Restart=always
RestartSec=5
Environment=GIN_MODE=release

# Environment file
EnvironmentFile=/root/tokogo/.env.production

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/tokogo/uploads

[Install]
WantedBy=multi-user.target
